# vps-agent/pennops-agent.service
[Unit]
Description=PennOps Autonomous Agent — Penn Enterprises LLC
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/agent
Environment=DEEPSEEK_API_KEY=PLACEHOLDER
ExecStart=/usr/bin/python3 /root/agent/agent-runner.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
